# ubuntu/service.conf
[supervisord]
nodaemon=true

[program:ssh]
command=/usr/sbin/service ssh start
autostart=true

[program:mysql]
command=/usr/sbin/service mysql start
autostart=true
